# docs/com.nexogrx.cowork-remote.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nexogrx.cowork-remote</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/YOUR_USERNAME/Projects/cowork-remote/server.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/YOUR_USERNAME/Projects/cowork-remote</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/YOUR_USERNAME/logs/cowork-remote.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USERNAME/logs/cowork-remote.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>4242</string>
    <key>COWORK_USER</key>
    <string>YOUR_USERNAME</string>
    <key>COWORK_PASS</key>
    <string>YOUR_PASSWORD</string>
    <key>COWORK_DIR</key>
    <string>/Users/YOUR_USERNAME/Projects</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
